# deploy/systemd/csgofriberg.service
[Unit]
Description=CSGO Friberg web game
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=csgofriberg
Group=csgofriberg
WorkingDirectory=/opt/csgofriberg
EnvironmentFile=/etc/csgofriberg/csgofriberg.env
UMask=0027

# A failed migration prevents the application from starting with a partial schema.
ExecStartPre=/usr/bin/node /opt/csgofriberg/server/dist/db/migrate.js
ExecStart=/usr/bin/node /opt/csgofriberg/server/dist/index.js

Restart=on-failure
RestartSec=5s
TimeoutStartSec=120s
TimeoutStopSec=30s
KillSignal=SIGTERM
SyslogIdentifier=csgofriberg

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
